# firmware/pi/miehome-firstboot.service
[Unit]
Description=MIEI HOME First Boot Setup
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/etc/mieihome/.firstboot_complete

[Service]
Type=oneshot
ExecStart=/opt/mieihome/first_boot.sh
RemainAfterExit=yes
StandardOutput=journal+console
StandardError=journal+console
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
